# deploy/voice-input-proxy.service
[Unit]
Description=Voice Input Auth Proxy
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/voice-input
Environment=VALID_API_KEY=sk-spense-93b10b0892fc4c1043aae0f74b1fe8a99b47585fda34e9b627ae43a4ca32a08c
ExecStart=/usr/bin/python3 auth-proxy.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
